# services/stg/compstat_notifier.timer
[Unit]
Description=compstat_notifier

[Timer]
OnCalendar=2015-05-20 *:*/5:00

[Install]
WantedBy=timers.target

[X-Fleet]
Global=true
