# services/stg/compstat_scraper.timer
[Unit]
Description=compstat_scraper

[Timer]
OnCalendar=Mon,Thu,Fri *-*-* 12:30

[Install]
WantedBy=timers.target

[X-Fleet]
Global=true
